# files/etc/cron.d/hello
#!/bin/bash
#

##---------------------
#!!CRONTAB	*/1 * * * *
##---------------------

function crontab_entry {
	INFO "hello world"
}
